# chip-input-reset/system/chip-local.service
[Unit]
ConditionFileIsExecutable=/usr/sbin/chip.local
Requires=lightdm.service
After=lightdm.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/chip.local
TimeoutSec=0
StandardOutput=tty
RemainAfterExit=yes
SysVStartPriority=99

[Install]
WantedBy=graphical.target
